# snort.list
/.
/etc
/etc/cron.daily
/etc/default
/etc/default/snort
/etc/init.d
/etc/init.d/snort
/etc/logrotate.d
/etc/logrotate.d/snort
/etc/ppp
/etc/ppp/ip-down.d
/etc/ppp/ip-down.d/snort
/etc/ppp/ip-up.d
/etc/ppp/ip-up.d/snort
/etc/snort
/usr
/usr/bin
/usr/bin/u2boat
/usr/bin/u2spewfoo
/usr/sbin
/usr/sbin/appid_detector_builder.sh
/usr/sbin/snort
/usr/sbin/u2boat
/usr/sbin/u2openappid
/usr/sbin/u2spewfoo
/usr/sbin/u2streamer
/usr/share
/usr/share/doc
/usr/share/doc/snort
/usr/share/doc/snort/AUTHORS
/usr/share/doc/snort/BUGS
/usr/share/doc/snort/CREDITS.gz
/usr/share/doc/snort/INSTALL.gz
/usr/share/doc/snort/NEWS.gz
/usr/share/doc/snort/README.Debian.gz
/usr/share/doc/snort/README.Maintainer
/usr/share/doc/snort/README.PLUGINS
/usr/share/doc/snort/README.csv
/usr/share/doc/snort/README.gz
/usr/share/doc/snort/RELEASE.NOTES
/usr/share/doc/snort/TODO.Debian.gz
/usr/share/doc/snort/USAGE.gz
/usr/share/doc/snort/copyright
/usr/share/doc/snort/generators
/usr/share/doc/snort/snort_rules.html
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/snort
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/snort.8.gz
/usr/share/man/man8/u2boat.8.gz
/usr/share/man/man8/u2spewfoo.8.gz
/var
/var/lib
/var/lib/snort
/var/log
/var/log/snort
/usr/share/doc/snort/NEWS.Debian.gz
/usr/share/doc/snort/changelog.Debian.gz
